# templates/etc/profile.d/st2.sh.erb
<%- if @_api_url -%>
ST2_API_URL=<%= @_api_url %>
<%- end -%>
<%- if @_auth_url -%>
ST2_AUTH_URL=<%= @_auth_url %>
<%- end -%>
