# templates/etc/profile.d/st2.sh.erb
## This file is maintained by Puppet

# URL to the StackStorm API service (st2api) 
<%- if @api_url -%>
export ST2_API_URL="<%= @api_url %>"
<%- else -%>
# export ST2_API_URL=""
<%- end -%>

# URL to the StackStorm auth service (st2auth) 
<%- if @auth_url -%>
export ST2_AUTH_URL="<%= @auth_url %>"
<%- else -%>
# export ST2_AUTH_URL=""
<%- end -%>
