# templates/etc/profile.d/st2.sh.erb
ST2_API_URL=<%= @api_url %>
ST2_AUTH_URL=<%= @auth_url %>

export ST2_API_URL ST2_AUTH_URL
