# templates/etc/systemd/system/mistral-api.service.erb
[Unit]
Description=Mistral Workflow API Service
BindsTo=mistral.service

[Service]
ExecStart=/opt/openstack/mistral/.venv/bin/gunicorn -b <%= @api_url %>:<%= @api_port %> -w 2 mistral.api.wsgi
WorkingDirectory=/opt/openstack/mistral
Restart=on-abort

[Install]
WantedBy=multi-user.target