# templates/etc/systemd/system/mistral.service.erb
[Unit]
Description=Mistral Workflow Service
Wants=mistral-api.service

[Service]
ExecStart=/opt/openstack/mistral/.venv/bin/python /opt/openstack/mistral/mistral/cmd/launch.py --config-file /etc/mistral/mistral.conf --log-file /var/log/mistral.log --log-config-append /etc/mistral/wf_trace_logging.conf --server "<%= @subsystems %>"
Restart=on-abort

[Install]
WantedBy=multi-user.target
