# templates/etc/systemd/system/st2service_multi.service.erb
[Unit]
Description=StackStorm service <% @name %>
After=network.target

[Service]
Type=simple
User=st2
Group=st2
Environment="DAEMON_ARGS=--config-file /etc/st2/st2.conf"
Environment="WORKERID=%i"
EnvironmentFile=-/etc/sysconfig/<%- @name %>
ExecStart=/usr/bin/<%- @name %> $DAEMON_ARGS
TimeoutSec=60
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
