# templates/etc/init/mistral-api.conf.erb
# Moving under puppet. To enable:
# 1. Install gunicorn to mistral virutal env /opt/stackstorm/mistral/.venv
# 2. Modify /etc/init/mistral.conf and remove "api" from --server
# 3. Modify st2ctl to add mistral-api
# 4. Remove these comments :)

description "Mistral Workflow API Service"

start on starting mistral
stop on stopping mistral

script
    chdir /opt/openstack/mistral
    exec /opt/openstack/mistral/.venv/bin/gunicorn -b <%= @api_url %>:<%= @api_port %> -w 2 mistral.api.wsgi
end script